# robominder.service
# Create a dedicated user
# sudo useradd -r -s /bin/false robominder

# Create directory and copy files
# sudo mkdir /opt/robominder
# sudo cp bot.py /opt/robominder/
# sudo chown -R robominder:robominder /opt/robominder

# Create service file
# sudo nano /etc/systemd/system/robominder.service
# (paste the service configuration from below)

# Enable and start the service
# sudo systemctl enable robominder
# sudo systemctl start robominder

# Check status
# sudo systemctl status robominder

# View logs
# journalctl -u robominder -f

[Unit]
Description=XMPP Reminder Bot
After=network.target

[Service]
Type=simple
User=robominder
Group=robominder
WorkingDirectory=/opt/robominder
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /opt/robominder/bot.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target